# usr/lib/systemd/user/rustdesk-user.service
[Unit]
Description=RustDesk User Service

[Service]
Type=simple
ExecStart=/usr/bin/rustdesk --service

[Install]
WantedBy=default.target